# scripts/deploy/speedero-snapshot.timer
[Unit]
Description=Run nightly snapshot at 00:30 PT

[Timer]
OnCalendar=*-*-* 00:30:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
