# scripts/deploy/speedero-snapshot-watchdog.service
[Unit]
Description=Speedero Security — capture watchdog
# No Requires=speedero-security.service: if the app is degraded the curl fails
# and the unit shows up in `systemctl --failed`, leaving a trace. The app-down
# class is out of scope (the watchdog lives inside the app it monitors).
After=network-online.target speedero-security.service

[Service]
Type=oneshot
User=andrew
EnvironmentFile=/data/SecApp/shared/.env.production
ExecStart=/usr/bin/curl -fsS \
  -H "Authorization: Bearer ${SNAPSHOT_CRON_TOKEN}" \
  -X POST http://127.0.0.1:3000/SecApp/api/snapshot/watchdog
StandardOutput=journal
StandardError=journal
