# scripts/deploy/speedero-security.service
[Unit]
Description=Speedero Security (Next.js, basePath /SecApp)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=andrew
Group=users
WorkingDirectory=/data/SecApp/current
EnvironmentFile=/data/SecApp/shared/.env.production
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NODE_OPTIONS=--max-old-space-size=1024
ExecStart=/usr/bin/node /data/SecApp/current/node_modules/next/dist/bin/next start -H 127.0.0.1 -p 3000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/data/SecApp

[Install]
WantedBy=multi-user.target
