# scripts/deploy/speedero-snapshot-watchdog.timer
[Unit]
Description=Check that SecApp snapshot capture is still running

[Timer]
# Fires 00:00 / 06:00 / 12:00 / 18:00 America/Los_Angeles.
# Verify: systemd-analyze calendar '*-*-* 00/6:00:00'
OnCalendar=*-*-* 00/6:00:00 America/Los_Angeles
Persistent=true

[Install]
WantedBy=timers.target
